# docker/local_supervisor/for_socialtech.conf
[program:for_superproject]
command=php /var/www/superproject/bin/console run-worker-async
process_name=%(program_name)s%(process_num)02d
user=www-data
numprocs=1
autostart=true
autorestart=true
startretries=10
stopsignal=INT
stopasgroup=true
stdout_logfile=/var/www/superproject/var/log/for_superproject_stdout.log
stderr_logfile=/var/www/superproject/var/log/for_superproject_error.log
